# com.djailla.empty/source/rc.local
#!/bin/sh -e
#
# rc.local
#

case $1 in
  start)
    echo 'start'
    ;;
  stop)
    echo 'stop'
    ;;
esac
